# dot_config/systemd/user/battery-notify.service
[Unit]
Documentation=https://github.com/cdown/battery-notify

[Service]
Type=notify
Environment=RUST_BACKTRACE=1
ExecStart=%h/bin/battery-notify
