# dot_config/systemd/user/habitual.service
[Service]
ExecStart=/usr/bin/bb %h/bin/habitual.clj
Restart=on-failure

[Install]
WantedBy=default.target
