# dot_config/systemd/user/habitual.service
[Service]
ExecStart=/usr/bin/bb %h/bin/habitual.clj -p "/usr/bin/notify-send -u critical"
Restart=on-failure

[Install]
WantedBy=default.target
