# dot_config/systemd/user/magdalena-cleanup.service
[Unit]
Description=Interactive shell history cleanup

[Service]
Type=oneshot
ExecStart=/opt/tools/magdalena --cleanup
