# dot_config/systemd/user/no-swap.slice
[Unit]
Description = No Swap but HVT for OOM
[Slice]
MemorySwapMax = 0
[Service]
OOMScoreAdjust = 1000
