# dot_config/systemd/user/nota.timer
[Unit]
Description=Run nota backup periodically

[Timer]
OnBootSec=30min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
