# dot_config/systemd/user/singlefile-sync.service
[Unit]
Description=SingleFile Sync Service
After=graphical-session.target

[Service]
Type=simple
ExecStart=/opt/tools/singlefile-sync
Restart=on-failure
RestartSec=30
KillMode=control-group
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
